# templates/jenkins-default.slice.erb
[Unit]
Description=Limited resources Slice
DefaultDependencies=no
Before=slices.target
[Slice]
CPUQuota=<%= @slice_percentage %>%
